# supervisord.conf
[supervisord]
nodaemon=true
loglevel=info
user=root

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:prometheus]
command=/opt/prometheus/prometheus --config.file=/etc/spark-monitor/prometheus.yml --storage.tsdb.path=/var/lib/prometheus --web.console.libraries=/opt/prometheus/console_libraries --web.console.templates=/opt/prometheus/consoles
autorestart=true
priority=100
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:grafana]
command=/opt/grafana/bin/grafana server --homepath=/opt/grafana
environment=GF_SECURITY_ADMIN_USER="admin",GF_SECURITY_ADMIN_PASSWORD="admin",GF_USERS_ALLOW_SIGN_UP="false",GF_PATHS_PROVISIONING="/etc/spark-monitor/grafana/provisioning"
autorestart=true
priority=200
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:node_exporter]
command=/opt/node_exporter/node_exporter --path.procfs=/host/proc --path.sysfs=/host/sys --collector.filesystem.ignored-mount-points=^/(dev|proc|sys|var/lib/docker/.+|var/lib/containers/.+)[:] --collector.filesystem.ignored-fs-types=^(autofs|binfmt_misc|bpf|cgroup2?|configfs|debugfs|devpts|devtmpfs|fusectl|hugetlbfs|iso9660|mqueue|nsfs|overlay|proc|procfs|pstore|rpc_pipefs|securityfs|selinuxfs|squashfs|sysfs|tmpfs)
autorestart=true
priority=200
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:gpu_exporter]
command=/opt/dcgm/dcgm-exporter --address=0.0.0.0:9400 -f /etc/spark-monitor/custom-dcgm-counters.csv
autorestart=true
priority=200
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:vllm_exporter]
command=python3 /app/exporters/vllm_exporter.py --port 8001 --scrape-interval 5
environment=VLLM_URL="http://host.docker.internal:8000"
autorestart=true
priority=300
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:system_exporter]
command=python3 /app/exporters/system_exporter.py --port 9106 --scrape-interval 15
autorestart=true
priority=300
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:agent_exporter]
command=python3 /app/exporters/agent_exporter.py --port 9107 --interval 60 --claude-stats /host/.claude/stats-cache.json --codex-db /host/.codex/state_5.sqlite --opencode-stats /host/.opencode/stats-cache.json --opencode-db /host/.opencode/state_5.sqlite
autorestart=true
priority=300
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
